# proximadb.service
[Unit]
Description=ProximaDB Vector Database Server
Documentation=https://docs.proximadb.ai
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=%i
Group=%i
ExecStart=/home/vsingh/code/proximadb/target/release/proximadb-server --config /home/vsingh/code/proximadb/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/data/proximadb

# Environment
Environment=RUST_LOG=info
Environment=PROXIMADB_CONFIG=/home/vsingh/code/proximadb/config.toml
Environment=PROXIMADB_DATA_DIR=/data/proximadb

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Working directory
WorkingDirectory=/home/vsingh/code/proximadb

# Standard output and error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proximadb

[Install]
WantedBy=default.target