# misc/services/proximadb-docker.service
[Unit]
Description=ProximaDB Vector Database Server (Docker Optimized)
Documentation=https://github.com/vjsingh1984/proximadb
After=network-online.target
Wants=network-online.target
StartLimitInterval=60s
StartLimitBurst=3

[Service]
Type=exec
User=proximadb
Group=proximadb
WorkingDirectory=/opt/proximadb
ExecStart=/opt/proximadb/bin/proximadb-server --config /opt/proximadb/config/proximadb.toml
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10

# Security hardening for production
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/proximadb/data /var/log/proximadb
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits optimized for vector database
LimitNOFILE=1048576
LimitNPROC=32768
LimitMEMLOCK=infinity
LimitCORE=infinity

# Memory and CPU management
MemoryAccounting=true
MemoryMax=80%
CPUAccounting=true
TasksMax=infinity

# Environment variables
Environment=RUST_LOG=info,proximadb=debug
Environment=RUST_BACKTRACE=1
Environment=PROXIMADB_CONFIG=/opt/proximadb/config/proximadb.toml
Environment=PROXIMADB_DATA_DIR=/opt/proximadb/data
Environment=PROXIMADB_LOG_DIR=/var/log/proximadb
Environment=PROXIMADB_BIND_ADDRESS=0.0.0.0:5678

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proximadb
SyslogFacility=daemon

# Health monitoring
WatchdogSec=60
NotifyAccess=main

[Install]
WantedBy=multi-user.target